# passwd_exp.cron.in
#!/bin/sh
#Run password expiration check program ( passwd_exp )
# by Samuel Behan <behan@frida.fri.utc.sk>
#

#setup command path replaced by install script
exec_prefix="@prefix@"
passwd_exp="@sbindir@/passwd_exp"

#at first run the passwd_exp program
RETVAL=`$passwd_exp $NO_DB 2>&1`

#check returned ( nothing is good ), on error sendmail to root
if [ -n "$RETVAL" ];
  then 
  DATE=`date`
  mail root -s "Error running passwd checks" <<__EOF__

  There occured some error(s) while running password
    expiration program checks by cron.

  Program command line:
$passwd_exp
  
  Program returned:
$RETVAL
  
---
  Mail generated by cron on $DATE.

__EOF__

fi

#free up enviroment
unset DATE
unset RETVAL
unset exec_prefix
unset passwd_exp

